# apple_scripts/com.workbench.lifx-dash-bulbs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.workbench.lifx-dash-bulbs</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/sh</string>
      <string>/u/code/workbench/bin/lifx_dash_bulbs</string>
    </array>
  </dict>
</plist>
